# docker/livefiles/liveroot/etc/rc.local
#!/bin/bash
echo "rc.local Script" 
echo "$(date) rc.local" >> /logf
#sleep 20
/root/startuprc.sh &
exit 0
